# files/root/etc/systemd/system/wg-quick@wg-mullvad.service.d/override.conf
[Service]
ExecStartPre=/usr/local/sbin/mullvad-select-relay
ExecStartPost=/usr/local/sbin/mullvad-open-port
ExecStopPost=/usr/local/sbin/mullvad-close-port