# stow/.local/etc/autostart/dunst.desktop
[Desktop Entry]
Name=dunst
Comment=a notification daemon
Exec=dunst
Type=Application
