# stow/.local/etc/autostart/picom.desktop
[Desktop Entry]
Name=picom
Comment=an X11 compositor
Exec=picom
Type=Application
