# stow/.local/etc/autostart/setbg.desktop
[Desktop Entry]
Name=setbg
Comment=Set the background image
Exec=setbg
Type=Application
