# stow/.local/etc/autostart/setbg.desktop
[Desktop Entry]
Name=setbg
Comment=restore background image
Exec=setbg
Type=Application
