# stow/.local/etc/autostart/sxhkd.desktop
[Desktop Entry]
Name=sxhkd
Comment=Define keyboard shortcuts
Exec=sxhkd
Type=Application
