# stow/.local/etc/autostart/sxhkd.desktop
[Desktop Entry]
Name=sxhkd
Comment=keyboard shortcut daemon
Exec=sxhkd
Type=Application
