# stow/.local/etc/autostart/udiskie.desktop
[Desktop Entry]
Name=udiskie
Comment=automounter for removable media
Exec=udiskie --automount --tray --notify
Type=Application
