# stow/.local/etc/autostart/unclutter.desktop
[Desktop Entry]
Name=unclutter
Comment=hide mouse cursor when idle
Exec=unclutter
Type=Application
