# stow/.local/etc/autostart/unclutter.desktop
[Desktop Entry]
Name=unclutter
Comment=Hide the mouse cursor while typing
Exec=unclutter
Type=Application
