# files/root/etc/systemd/system/wg-quick@mullvad.service.d/override.conf
[Service]
ExecStartPre=/usr/local/sbin/mullvad-select-relay
