# stow/.local/etc/systemd/user/notify-battery.timer
[Unit]
Description=check discharging batteries periodically

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
