# stow/.local/etc/systemd/user/notify-battery.service
[Unit]
Description=notify about discharging batteries

[Service]
ExecStart=%h/.local/bin/notify-battery

