# stow/.local/etc/systemd/user/sway-setbg.service
[Unit]
BindsTo=sway-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/sway-setbg

[Install]
WantedBy=sway-session.target
