# stow/.local/etc/systemd/user/udiskie.service
[Unit]
BindsTo=sway-session.target

[Service]
Type=simple
ExecStart=/usr/bin/udiskie --automount --notify

[Install]
WantedBy=sway-session.target
