# stow/.config/systemd/user/check-battery.service
[Unit]
Description=check battery power

[Service]
ExecStart=%h/.local/bin/check-battery
