# stow/.config/autostart/dunst.desktop
[Desktop Entry]
Name=dunst
Comment=A notification daemon
Exec=dunst
Type=Application
