# stow/.config/autostart/xrandr-primary.desktop
[Desktop Entry]
Name=xrandr-primary
Comment=Set the
Exec=unclutter
Type=Application
