# files/home/.config/systemd/user/check-battery.service
[Service]
ExecStart=%h/.local/bin/check-battery
