# files/home/.config/systemd/user/check-battery.service
[Unit]
Description=Check battery power

[Service]
ExecStart=%h/.local/bin/check-battery
