# home/.config/systemd/user/check-battery.timer
[Unit]
Description=periodically check battery power

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
